# scripts/ec2/services/homeguard-trading.target
[Unit]
Description=Homeguard Trading Services
After=network.target

[Install]
WantedBy=multi-user.target
